# web/supervisord_web.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
loglevel=info
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
prompt=supervisorctl

[program:gunicorn]
command=gunicorn colander2:app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

;[program:celery_beat]
;command=/bin/bash -c "exec celery --app=app:celery_app beat -S redbeat.RedBeatScheduler --loglevel=INFO"
;directory=/web
;numprocs=1
;stdout_logfile=/dev/fd/1
;stdout_logfile_maxbytes=0
;redirect_stderr=true
;autostart=true
;autorestart=true
;startsecs=10
;stopwaitsecs=60
;stopasgroup=true
;priority=1000
